# scripts/mission_runner.service
[Unit]
Description=Mission Runner Service
Documentation=https://github.com/your-repo/sih-project
After=network.target mosquitto.service telemetry.service
Wants=mosquitto.service telemetry.service
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/sih
Environment=PYTHONPATH=/home/pi/sih/src
Environment=DEVICE_ID=pi-drone-01
Environment=MQTT_HOST=localhost
Environment=MQTT_PORT=1883
Environment=MISSION_SPEED=5.0
Environment=MISSION_ALTITUDE=50.0

# Load MQTT credentials if available
EnvironmentFile=-/home/pi/.env.mqtt

# Python virtual environment
ExecStart=/home/pi/sih/.venv/bin/python /home/pi/sih/src/mission_runner.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/drone/logs

# Resource limits
LimitNOFILE=1024
MemoryMax=256M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mission-runner

[Install]
WantedBy=multi-user.target
